# deploy/d-brain-process.service
[Unit]
Description=d-brain Daily Processing

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/projects/agent-second-brain
ExecStart=/home/ubuntu/projects/agent-second-brain/scripts/process.sh
Environment=PYTHONUNBUFFERED=1
